# bin/crontabs/work-laptop
# m h  dom mon dow   command
* * * * * cd /Users/bernard/fieldbook && javascript-ctags 'client/**/*.js' 'server/**/*.js' 'common/**/*.js'
