# bin/mac/clipShare.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>clipShare</string>
    <key>Program</key>
    <string>/Users/bernard/bin/mac/screen-copy-sync.sh</string>
    <key>WatchPaths</key>
    <array>
      <string>/tmp/screen-exchange</string>
    </array>
    <key>ThrottleInterval</key>
    <integer>1</integer>
  </dict>
</plist>
